# ingest/crontab
0 * * * * python /main/ingest/fastforward.py >/dev/null 2>&1
